# robot_drive.service
[Unit]
Description=Robot
After=xdk-daemon.service

[Service]
#User=root
ExecStart=/usr/bin/http-server /home/root/drive-interface

[Install]
WantedBy=multi-user.target
